# cyt.service
[Unit]
Description=Chasing Your Tail - Wireless Surveillance Detection
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/media/psf/my_projects/0_active_projects/Chasing-Your-Tail-NG
ExecStart=/usr/bin/python3 cyt_daemon.py start
ExecStop=/usr/bin/python3 cyt_daemon.py stop
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
